# backend/cmd/stoarama-relay/templates/systemd.service.tmpl
[Unit]
Description=Stoarama Relay
After=network-online.target
Wants=network-online.target

[Service]
ExecStart={{.ExePath}} run
Restart=on-failure
RestartSec=15s

[Install]
WantedBy=default.target
